# src/99-detectar-usb.rules
SUBSYSTEM=="block", KERNEL=="sd[a-z]", ACTION=="add", RUN+="/usr/local/bin/usb_detectada.sh"
